# crontab
*/1 * * * * cd code/qbtmindf && ./cronjob
